# ospool-ep-container/gratia-probe-ospool-ep-container.cron
0,15,30,45 * * * * root /usr/share/gratia/common/cron_check  /etc/gratia/ospool-ep-container/ProbeConfig && /usr/share/gratia/ospool-ep-container/osgpilot_meter -s 900
